# etc/systemd/check-ttymidi-usage.timer
[Unit]
Description=Check for ttymidi CPU Usage every minute

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
