# etc/systemd/rfkill-unblock-all.service
[Unit]
Description=Unblock all wireless devices

[Service]
Type=oneshot
ExecStart=/usr/sbin/rfkill unblock all

[Install]
WantedBy=multi-user.target
